# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:loftgw]
command=/app/loftgw
directory=/root
# 【调整】：不再对 Loft-Gateway 写死过紧的限制。
# 移除 GOMEMLIMIT 的硬性压制，仅保留 MALLOC_ARENA_MAX=2 防止底层 Cgo 的内存碎片膨胀。
# 将内存的管理权下放给 PHP 控制面下发的 sys_mem_limit_mb 参数，进行动态 Go 内存软限制。
environment=MALLOC_ARENA_MAX="2"
autostart=true
autorestart=true
# 增强稳定性：程序跑满 5 秒才算启动成功
startsecs=5
# 增强稳定性：失败后最多重试 10 次（默认是 3 次）
startretries=10
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:cloudflared]
# 使用 http2 协议替代 quic，大幅降低小内存容器的内存占用
command=/app/cloudflared tunnel --no-autoupdate --protocol http2 run --token %(ENV_tunnelid)s
directory=/app
# 【调整】：适度收紧 Cloudflared 的内存限制
# 1. GOMEMLIMIT="95MiB"：在 128MB 的基础上进一步压制（对于纯转发隧道通常足够）。
# 2. GOGC="30"：让垃圾回收更积极（之前为 50），代价是极其微小的 CPU 占用增加，但能有效防止内存水位随连接数堆叠。
# 3. MALLOC_ARENA_MAX="1"：强制所有线程共享单一内存池，最大限度降低底层的内存预分配开销。
environment=GOMEMLIMIT="95MiB",GOGC="30",MALLOC_ARENA_MAX="1"
autostart=true
autorestart=true
startsecs=5
startretries=10
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
