# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:loftgw]
command=/app/loftgw
directory=/root
# 优化点：
# 1. GOGC=30：让 Xray/Mihomo 的 GC 更加频繁，消耗少许 CPU 降低常驻内存。
# 2. GODEBUG=madvdontneed=1：强制立刻将释放的内存归还给 Linux 内核，避免虚高。
# 3. MALLOC_ARENA_MAX=1：彻底限制 CGO/glibc 层的内存分配区，极致降低碎片。
environment=GOMEMLIMIT="256MiB",GOGC="30",GODEBUG="madvdontneed=1",MALLOC_ARENA_MAX="1"
autostart=true
autorestart=true
startsecs=5
startretries=10
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:cloudflared]
# 维持 http2 的低内存连接池
command=/app/cloudflared tunnel --no-autoupdate --protocol http2 run --token %(ENV_tunnelid)s
directory=/app
# 优化点：
# 1. GOGC=20：对 Cloudflared 执行最激进的内存回收（只要对象不用了，立刻烧 CPU 回收）。
# 2. GOMEMLIMIT="96MiB"（或比您平时观察到的稍高一点）：一旦触碰，强行用 CPU 压平。
# 3. 同样加上 madvdontneed=1
environment=GOMEMLIMIT="96MiB",GOGC="20",GODEBUG="madvdontneed=1",MALLOC_ARENA_MAX="1"
autostart=true
autorestart=true
startsecs=5
startretries=10
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
